# srcs/mysql/supervisor.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=mysqld --user=root --datadir=/var/lib/mysql
