# srcs/ftps/supervisor.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=vsftpd /etc/vsftpd.conf

[program:telegraf]
command=/telegraf/telegraf
